# setups/systemd/electricsheep-hd-client.service.template
[Unit]
Description=Renders sheeps in HD

[Service]
ExecStart=$DAEMON_PATH
User=$USER
Restart=on-abort

[Install]
WantedBy=multi-user.target
